# root/etc/services.d/mylar/run
#!/usr/bin/with-contenv bash

exec \
	pip3 install -r requirements.txt

exec \
	s6-setuidgid abc python3 /app/mylar/Mylar.py --quiet --nolaunch \
	--datadir /config/mylar
